# tools/crontab/finna-metadata-collect
# Adjust paths as needed, deploy with "crontab finna-metadata-collect"

@monthly /path/to/finna-metadata-collect/collect.sh
